# src/etc/systemd/system/sim7600_4G_hat_init.timer
[Unit]
Description=Timer for the Waveshare SIM7600 4G Hat initialization script

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target